# com.breeze-mcp.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.breeze-mcp.server</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ivan/github/casualjim/breeze/.venv/bin/python</string>
        <string>-m</string>
        <string>breeze</string>
        <string>serve</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/ivan/github/casualjim/breeze</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ivan/github/casualjim/breeze/.venv/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>BREEZE_HOST</key>
        <string>127.0.0.1</string>
        <key>BREEZE_PORT</key>
        <string>9483</string>
        <key>BREEZE_DATA_ROOT</key>
        <string>/Users/ivan/Library/Application Support/breeze</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>VIRTUAL_ENV</key>
        <string>/Users/ivan/github/casualjim/breeze/.venv</string>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/breeze-mcp.log</string>
    
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/breeze-mcp.error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>